# release/linux/fpm.before-remove.sh
systemctl stop walletd.service
systemctl disable walletd.service
